# i2c_python_stats_stop.timer
[Unit]
Description="Schedulable i2c monitor for stop service"

[Timer]
OnCalendar=*-*-* 23:30:00
Persistent=true

[Install]
WantedBy=timers.target
