# systemd/raspberrypi_stats_cpp.service
[Unit]
Description=Raspberry Pi OLED Stats (C++ low-memory service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/raspberrypi_stats_cpp
Restart=always
RestartSec=5
User=pi
Group=pi
EnvironmentFile=-/etc/default/raspberrypi_stats_cpp
WorkingDirectory=/usr/local/bin
StandardOutput=journal
StandardError=journal
Nice=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
